# backend/netgarde-domain-classifier.service
[Unit]
Description=NetGarde Domain Classification Worker
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/ubuntu/netgarde/backend
EnvironmentFile=/home/ubuntu/netgarde/backend/.env.production
ExecStart=/usr/bin/python3 -m app.workers.domain_classification_worker
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target

