# dns-sync/netgarde-blocked-sites-listener.service
[Unit]
Description=NetGarde DNS Policy DB Listener (LISTEN/NOTIFY -> run-sync.sh)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/ubuntu/netgarde
EnvironmentFile=/etc/netgarde/backend.env
Environment=DB_NOTIFY_CHANNELS=blocked_sites_changed,policy_changed
Environment=SYNC_COMMAND=/home/ubuntu/netgarde/dns-sync/run-sync.sh
ExecStart=/usr/bin/python3 /home/ubuntu/netgarde/dns-sync/blocked_sites_listener.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
