# dns-sync/trustedge-log-watcher.service
[Unit]
Description=TrustEdge DNS Log Watcher - Real-time DNS query streaming
After=network.target docker.service
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/ubuntu/trustedge/dns-sync/dns_log_watcher.py
Restart=always
RestartSec=5
User=root
WorkingDirectory=/home/ubuntu/trustedge/dns-sync

# Environment variables (override defaults as needed)
Environment=DNSMASQ_LOG_PATH=/var/log/dnsmasq.log
Environment=STATE_FILE_PATH=/var/lib/trustedge/log_parser_state
Environment=BLOCKED_DOMAINS_PATH=/etc/dnsmasq.d/blocked-domains.conf
Environment=API_BASE_URL=http://localhost:8000
Environment=BATCH_SIZE=50
Environment=POLL_INTERVAL=2
Environment=FLUSH_INTERVAL=3
Environment=LOG_JSON=1
Environment=LOG_LEVEL=INFO
# Zero Trust: require this when backend DNS_INGEST_TOKEN is set
# Environment=DNS_INGEST_TOKEN=REPLACE_WITH_LONG_RANDOM_SECRET

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trustedge-log-watcher

[Install]
WantedBy=multi-user.target

